# packages/apply-planner/test/golden/crontab-example-code-projects.cron
# GENERATED BY synccenter — do not edit
# SchedulePlan → crontab fragment

# example-code-projects → qnap-ts453d
*/15 * * * * docker exec rclone-rcd rclone bisync /share/Sync/code gdrive:sync/code --filters-file=/share/synccenter-config/compiled/dev-monorepo/filter.rclone --resilient --recover --max-lock=2m --conflict-resolve=newer --conflict-loser=pathrename
